# vagrant/websockets.conf
[program:websockets]
command=php /home/vagrant/code/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=vagrant